# docs/2012/root.crontab
# minute hour day_of_month month day_of_week command
# Reboot 10:05 every Monday morning
5 10 * * Mon /sbin/reboot
